# locale/profile.d/boloot-calendar.sh
# BOLOOT Persian Calendar — boloot.ir
# Installed to /etc/profile.d/boloot-calendar.sh by boloot-calendar-ctl apply-system

export LC_TIME="${BOLOOT_LC_TIME:-fa_IR@calendar=persian}"
export BOLOOT_CALENDAR_PROFILE="${BOLOOT_CALENDAR_PROFILE:-iran}"
export BOLOOT_CALENDAR_TIMEZONE="${BOLOOT_CALENDAR_TIMEZONE:-Asia/Tehran}"
